# .config/systemd/user/waybar.service
[Unit]
Description=Waybar status bar
After=graphical-session.target
BindsTo=graphical-session.target
Requires=graphical-session.target
Wants=graphical-session.target  

[Service]
Type=simple
Environment=WAYLAND_DISPLAY=%t/wayland-1
Environment=XDG_RUNTIME_DIR=%t
Environment=DISPLAY=:0  # For XWayland compatibility
ExecStart=/usr/bin/waybar
Restart=on-failure
RestartSec=3
KillMode=mixed
TimeoutStopSec=2

[Install]
WantedBy=graphical-session.target

